# utilities/lstu-minion@.service
[Unit]
Description=Shortened URLs service job queue
Documentation=https://framagit.org/fiat-tux/hat-softwares/lstu

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/lstu/
ExecStart=/usr/local/bin/carton exec script/application minion worker

[Install]
WantedBy=multi-user.target
